# software/_deploy/nx/setup/jetson_clock_service/jetsonClocks.service
[Unit]
Description=Maximize Jetson Performance
After=multi-user.target

[Service]
ExecStart=/usr/bin/jetson_clocks --fan

[Install]
WantedBy=multi-user.target
